# fhem/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:fhem]
command=perl /opt/fhem/fhem.pl configDB
